# root/etc/services.d/rclone/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /usr/bin/rclone --config /config/.rclone.conf mount $RcloneRemotePath /mnt/rclone $RcloneMountCommands
