# root/etc/services.d/unionfs/run
#!/usr/bin/with-contenv bash

while [ ! -f $FILE_TO_WAIT ]
do
  echo "MergerFS mount not yet active"
  sleep 2
done

echo "MergerFS mount active"

exec s6-setuidgid abc /usr/bin/mergerfs -f -o defaults,nonempty,category.action=all,category.create=epff,minfreespace=0,allow_other,auto_cache,sync_read,umask=0002,func.getattr=newest /mnt/local=RW:/mnt/rclone=RW /mnt/unionfs

